# supervisord.conf
[program:goonpud]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py goonpugd
